# icy-pi-client/debian/icy-pi-client.service
[Unit]
Description=icy-pi-client service

[Service]
ExecStart=/usr/local/bin/icy-pi-client
Restart=always

[Install]
WantedBy=multi-user.target
